# etc/profile.d/xdg-cache-cuda.sh
export XDG_CACHE_HOME="${XDG_CACHE_HOME:-$HOME/.cache}"
export CUDA_CACHE_PATH="$XDG_CACHE_HOME/nv"

